# matviews/crontab
0 20 * * * /home/kumykov/refresh_data.sh
