# matviews/crontab
* * * * * /home/kumykov/refresh_data.sh
